# rootfs/etc/services.d/observer/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Supervisor observer service
# ==============================================================================

exec observer 
